# .config/systemd/user/dot-sync.service
[Unit]
Description=Auto Git Sync Service
After=network.target

[Service]
Type=simple
WorkingDirectory=%h
ExecStart=%h/dot-sync
