# .config/systemd/user/dot-sync.timer
[Unit]
Description=Run Git Sync Every 10 Minutes

[Timer]
OnBootSec=10min
OnUnitActiveSec=10min
Persistent=true

[Install]
WantedBy=timers.target
